# vibehub.service
[Unit]
Description=Vibehub - IDE for Claude Sessions
After=network.target

[Service]
Type=simple
ExecStart=/opt/vibehub/start.sh
Environment=NODE_ENV=production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
