# vibehub.service
[Unit]
Description=Vibehub - IDE for Claude Sessions
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/root/.env_secrets
WorkingDirectory=/root/projects/vibehub
ExecStart=/bin/bash -c 'set -a; source /root/.env_secrets; node build/main.js --conf conf/config.json5'
Environment=NODE_ENV=production
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
